# root/etc/services.d/tvheadend/run
#!/usr/bin/with-contenv bash

IFS=" " read -r -a RUN_ARRAY <<< "$RUN_OPTS"

exec \
        nice -n -15 ionice -c 1 -n 0 s6-setuidgid abc tvheadend -C -c /config "${RUN_ARRAY[@]}"
